# Desktop/autostart/custom_imwheel.desktop
[Desktop Entry]
Type=Application
Name=custom_imwheel
NoDisplay=false
Hidden=false
Exec=/usr/bin/imwheel -k
X-GNOME-Autostart-enabled=true
