# Desktop/autostart/startup.imwheel.desktop
[Desktop Entry]
Type=Application
Name=startup/imwheel
NoDisplay=false
Hidden=false
Exec=/usr/bin/imwheel -k
X-GNOME-Autostart-enabled=true
