# link_Desktop/crontab/in_move_ff_ext.sh
#!/bin/bash
# Extract the folder names from the full paths
if [[ $1 =~ ^.*[0-9]{4}-[0-9]{2}-[0-9]{2}.*$ ]]; then
  	mv "/home/lwh/Downloads/$1" "/home/lwh/Downloads/Download_Home/firefox/"
fi
done

