# link_Desktop/crontab/lwh_upgrade.sh
#!/bin/bash
#fupdate
#flatpak update --noninteractive --assumeyes
#flatpak uninstall --unused --noninteractive --assumeyes

#rm -rf /home/lwh/.var/app/com.opera.Opera/config/opera/History

#gupdate
#/home/lwh/.local/bin/gupdate > /dev/null 2>&1

#dupdate
#/uddsr/bin/distrobox-upgrade --all > /dev/null 2>&1
