# link_Desktop/autostart/bluetooth.desktop
[Desktop Entry]
Hidden=false
Name=bluetooth
Comment=
Terminal=false
Exec=/usr/bin/bluetoothctl connect 74:45:CE:6A:C3:FF
Type=Application
