# link_Desktop/autostart/conky.desktop
[Desktop Entry]
Hidden=false
Name=conky
Comment=
Terminal=false
Exec=conky --pause=10
Type=Application
