# link_Desktop/autostart/USB.desktop
[Desktop Entry]
Hidden=false
Name=USB
Comment=
Terminal=true
Exec=/home/lwh/.config/autostart/USB.ignition_delay.sh
Type=Application
