# script/crontab/auto_update.sh
#!/bin/bash

flatpak update --noninteractive --assumeyes

# cat ~/.key | sudo -S sudo apt update && sudo apt -y upgrade 